# production/carton_line/deploy/systemd/visionops-hp60c-sdk-bridge-watchdog.timer
[Unit]
Description=Run VisionOps HP60C Camera Recovery Watchdog

[Timer]
OnBootSec=25s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
