# production/carton_line/deploy/systemd/visionops-v3-tcp-pick.service
[Unit]
Description=VisionOps v3 Carton Line - Tube Pick TCP Client
After=network-online.target visionops-v3-runtime-pick.service
Wants=network-online.target visionops-v3-runtime-pick.service

[Service]
Type=simple
User=neardi
WorkingDirectory=/opt/visionops_v3
EnvironmentFile=-/etc/visionops_v3/carton_line.env
ExecStart=/opt/visionops_v3/production/carton_line/scripts/start_tcp_pick.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
