# production/carton_palletizing/deploy/systemd/visionops-v3-carton-palletizing-runtime.service
[Unit]
Description=VisionOps v3 Carton Palletizing RKNN Runtime
After=network-online.target visionops-orbbec336l-bridge.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/visionops_v3
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/visionops_v3/carton_palletizing.env
ExecStart=/opt/visionops_v3/production/carton_palletizing/scripts/start_runtime.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
